# conf/supervisor_elasticsearch.conf
[program:elasticsearch]
command=/opt/elasticsearch-1.4.4/bin/elasticsearch
